# software/console/systemd/greenhouse-switches.service
[Unit]
Description=Greenhouse Switches
Wants=init-pins.service
After=init-pins.service

[Service]
Type=simple
RemainAfterExit=no
ExecStart=/opt/greenhouse/bin/switches-mqtt-pub

[Install]
WantedBy=multi-user.target
